# overlays/debug/etc/systemd/getty@tty2.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/usr/sbin/agetty -o '-p -f -- \\u' --autologin core --noclear %I $TERM
